# deploy/247-chenzhu.conf
[program:flask_gunicorn]
directory=/workspace/power_model_system
environment=FLASK_CONFIG='chenzhu247',PYTHONUNBUFFERED=1
command=python3 ./run_flask.py
user=root
autorestart=true
startsecs=3
startretires=3
redirect_stderr=true
stopsignal=HUP
killasgroup=true
stopasgroup=true
stdout_logfile_maxbytes=100MB
stdout_logfile_backups = 20
stdout_logfile=/var/log/supervisor/flask_gunicorn.log
# stderr_logfile=/var/log/supervisor/flask_gunicorn.err
# stderr_logfile_maxbytes=1MB
# stderr_logfile_backups=10


[program:socketio_backend]
directory=/workspace/power_model_system/
environment=PYTHONPATH="/workspace/power_model_system/:/workspace/power_model_system/flask_app",PYTHONUNBUFFERED=1
command=python3 ./flask_app/websocket/socketio_backend.py
user=root
autorestart=true
startsecs=3
startretires=3
redirect_stderr=true
stdout_logfile_maxbytes=100MB
stdout_logfile_backups = 20
stdout_logfile=/var/log/supervisor/socketio_backend.log
# stderr_logfile=/var/log/supervisor/socketio_backend.err
# stderr_logfile_maxbytes=1MB
# stderr_logfile_backups=10

[program:apscheduler]
directory=/workspace/power_model_system/
environment=PYTHONPATH="/workspace/power_model_system/:/workspace/power_model_system/flask_app",PYTHONUNBUFFERED=1
command=python3 ./flask_app/apschedule_tasks/aps_main.py
user=root
autorestart=true
startsecs=3
startretires=3
redirect_stderr=true
stdout_logfile_maxbytes=100MB
stdout_logfile_backups = 20
stdout_logfile=/var/log/supervisor/apscheduler.log
# stderr_logfile=/var/log/supervisor/socketio_backend.err
# stderr_logfile_maxbytes=1MB
# stderr_logfile_backups=10