# deploy/251-yww-dev.conf
[program:celery]
directory=/workspace/power_model_system
environment=PYTHONPATH="/workspace/power_model_system/:/workspace/power_model_system/flask_app",PYTHONUNBUFFERED=1
command=/bin/bash -c "python3 ./flask_app/celery_tasks/main.py"
user=root
autorestart=true
startsecs=2
startretires=3
redirect_stderr=true
stdout_logfile_maxbytes=100MB
stdout_logfile_backups = 20
stdout_logfile=/var/log/supervisor/celery.log
# stderr_logfile=/var/log/supervisor/celery.err
# stderr_logfile_maxbytes=1MB
# stderr_logfile_backups=10

[program:fake_data]
directory=/workspace/power_model_system
environment=PYTHONPATH="/workspace/power_model_system/:/workspace/power_model_system/flask_app",PYTHONUNBUFFERED=1
command=python3 ./flask_app/util/common/fake_data.py
user=root
startsecs=5
startretires=3
redirect_stderr=true
stdout_logfile_maxbytes=100MB
stdout_logfile_backups = 20
stdout_logfile=/var/log/supervisor/fake_data.log
# stderr_logfile=/var/log/supervisor/fake_data.err
# stderr_logfile_maxbytes=1MB
# stderr_logfile_backups=10

[program:flask_gunicorn]
directory=/workspace/power_model_system
environment=FLASK_CONFIG='development',PYTHONUNBUFFERED=1,LD_LIBRARY_PATH=/usr/local/nvidia/lib:/usr/local/nvidia/lib64:/usr/local/cuda-11.7/targets/x86_64-linux/lib
command=python3 ./run_flask.py
user=root
autorestart=true
startsecs=3
startretires=3
redirect_stderr=true
stopsignal=HUP
killasgroup=true
stopasgroup=true
stdout_logfile_maxbytes=100MB
stdout_logfile_backups = 20
stdout_logfile=/var/log/supervisor/flask_gunicorn.log
# stderr_logfile=/var/log/supervisor/flask_gunicorn.err
# stderr_logfile_maxbytes=1MB
# stderr_logfile_backups=10


[program:socketio_backend]
directory=/workspace/power_model_system/flask_app/websocket
environment=PYTHONPATH="/workspace/power_model_system/:/workspace/power_model_system/flask_app",PYTHONUNBUFFERED=1
command=gunicorn -k eventlet -w 1 --bind 0.0.0.0:18010  socketio_backend:flask_app
user=root
autorestart=true
startsecs=3
startretires=3
redirect_stderr=true
stdout_logfile_maxbytes=100MB
stdout_logfile_backups = 20
stdout_logfile=/var/log/supervisor/socketio_backend.log
# stderr_logfile=/var/log/supervisor/socketio_backend.err
# stderr_logfile_maxbytes=1MB
# stderr_logfile_backups=10



[program:apscheduler]
directory=/workspace/power_model_system/
environment=PYTHONPATH="/workspace/power_model_system/:/workspace/power_model_system/flask_app",PYTHONUNBUFFERED=1
command=python3 ./flask_app/apschedule_tasks/aps_main.py
user=root
autorestart=true
startsecs=3
startretires=3
redirect_stderr=true
stdout_logfile_maxbytes=100MB
stdout_logfile_backups = 20
stdout_logfile=/var/log/supervisor/apscheduler.log
# stderr_logfile=/var/log/supervisor/socketio_backend.err
# stderr_logfile_maxbytes=1MB
# stderr_logfile_backups=10
2023-07-11 15:27:57.402 | INFO     | flask_app.apschedule_tasks.point_check.point_check:point_check_task:224 - {'batch': 6, 'point_judge_cost': 0.08189940452575684}
2023-07-11 15:27:57.558 | INFO     | flask_app.apschedule_tasks.point_check.point_check:oper_guide_task:96 - oper guide task cost : 0.35028576850891113