# deploy/252-production.conf
[program:celery]
directory=/workspace/power_model_system
command=/bin/bash -c "source /root/.virtualenvs/celery_virtual/bin/activate && export LD_LIBRARY_PATH=/usr/local/nvidia/lib:/usr/local/nvidia/lib64:/usr/local/cuda-11.7/targets/x86_64-linux/lib && celery -A flask_app.celery_tasks.tasks worker -l info -c 400 -n worker0 --heartbeat-interval=30"
environment=FLASK_CONFIG='dandong252',LD_LIBRARY_PATH=/usr/local/nvidia/lib:/usr/local/nvidia/lib64:/usr/local/cuda-11.7/targets/x86_64-linux/lib
user=root
autorestart=true
startsecs=2
startretires=3
redirect_stderr=true
killasgroup=true
stopasgroup=true
stdout_logfile_maxbytes=100MB
stdout_logfile_backups = 20
stdout_logfile=/var/log/supervisor/celery.log
# stderr_logfile=/var/log/supervisor/celery.err
# stderr_logfile_maxbytes=1MB
# stderr_logfile_backups=10


[program:flask_gunicorn]
directory=/workspace/power_model_system
environment=FLASK_CONFIG='dandong252',PYTHONUNBUFFERED=1,LD_LIBRARY_PATH=/usr/local/nvidia/lib:/usr/local/nvidia/lib64:/usr/local/cuda-11.7/targets/x86_64-linux/lib
command=python3 ./run_flask.py
user=root
autorestart=true
startsecs=3
startretires=3
redirect_stderr=true
killasgroup=true
stopasgroup=true
stdout_logfile_maxbytes=100MB
stdout_logfile_backups = 20
stdout_logfile=/var/log/supervisor/flask_gunicorn.log
# stderr_logfile=/var/log/supervisor/flask_gunicorn.err
# stderr_logfile_maxbytes=1MB
# stderr_logfile_backups=10


[program:socketio_backend]
directory=/workspace/power_model_system/flask_app/websocket
environment=FLASK_CONFIG='dandong252',PYTHONPATH="/workspace/power_model_system/:/workspace/power_model_system/flask_app",PYTHONUNBUFFERED=1
command=gunicorn -k eventlet -w 1 --bind 0.0.0.0:18010  socketio_backend:flask_app
user=root
autorestart=true
startsecs=3
startretires=3
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile_maxbytes=100MB
stdout_logfile_backups = 20
stdout_logfile=/var/log/supervisor/socketio_backend.log
# stderr_logfile=/var/log/supervisor/socketio_backend.err
# stderr_logfile_maxbytes=1MB
# stderr_logfile_backups=10



[program:apscheduler]
directory=/workspace/power_model_system/
environment=FLASK_CONFIG='dandong252',LD_LIBRARY_PATH=/usr/local/nvidia/lib:/usr/local/nvidia/lib64:/usr/local/cuda-11.7/targets/x86_64-linux/lib,PYTHONPATH="/workspace/power_model_system/:/workspace/power_model_system/flask_app",PYTHONUNBUFFERED=1
command=python3 ./flask_app/apschedule_tasks/aps_main.py
user=root
autorestart=true
startsecs=3
killasgroup=true
stopasgroup=true
startretires=3
redirect_stderr=true
stdout_logfile_maxbytes=100MB
stdout_logfile_backups = 20
stdout_logfile=/var/log/supervisor/apscheduler.log
# stderr_logfile=/var/log/supervisor/socketio_backend.err
# stderr_logfile_maxbytes=1MB
# stderr_logfile_backups=10
